# src/rsync_backup.timer
[Unit]
Description=Rsync Backup Service for timer

[Timer]
OnCalendar=Mon-Sun *:0:00
Persistent=true

[Install]
WantedBy=timers.target