# src/rsync_backup.service
[Unit]
Description=Rsync Backup Service for main

[Service]
User=root
ExecStart=/usr/local/bin/rsync_backup.sh

[Install]
WantedBy=timers.target